# files/etc/crontabs/job_script.sh
#!/bin/sh

# Code goes here.
# echo "This is a script, run by cron!"