# files/etc/service/cron/run
#!/bin/bash

# Conf. job and start cron
/usr/bin/crontab /etc/crontabs/job_conf && /usr/sbin/crond -b -l 8